# server/scraper/com.lofi-atc.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lofi-atc.scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hojin/repos/lofi-atc/server/scraper/.venv/bin/python</string>
        <string>/Users/hojin/repos/lofi-atc/server/scraper/top_feeds_playwright.py</string>
        <string>--daemon</string>
        <string>--interval</string>
        <string>30</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hojin/repos/lofi-atc/server/scraper</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/hojin/repos/lofi-atc/server/scraper/scraper.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hojin/repos/lofi-atc/server/scraper/scraper.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/hojin/repos/lofi-atc/server/scraper/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
